# systemd/dye.socket
[Socket]
ListenFIFO=%t/dye.stdin
Service=dye.service
